# templates/network-ringbuffer.service.tmpl
[Unit]
Description=Network Ring Buffer Tuning
After=network.target
Wants=network.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/network-ringbuffer.sh

[Install]
WantedBy=multi-user.target
