# templates/stealth-firewall.service.tmpl
[Unit]
Description=Stealth Firewall - Block all incoming connections on public IP
After=network-online.target
Wants=network-online.target
Before=ssh.service

[Service]
Type=oneshot
RemainAfterExit=yes

# Flush existing rules
ExecStartPre=/sbin/iptables -F INPUT
ExecStartPre=/sbin/iptables -F OUTPUT

# Set default policies
ExecStartPre=/sbin/iptables -P INPUT DROP
ExecStartPre=/sbin/iptables -P FORWARD ACCEPT
ExecStartPre=/sbin/iptables -P OUTPUT ACCEPT

# Allow loopback
ExecStartPre=/sbin/iptables -A INPUT -i lo -j ACCEPT

# Allow established and related connections (for outgoing internet)
ExecStartPre=/sbin/iptables -A INPUT -m state --state ESTABLISHED,RELATED -j ACCEPT

# Allow all traffic on internal bridge (vmbr0/vmbr1)
ExecStartPre=/sbin/iptables -A INPUT -i vmbr0 -j ACCEPT
ExecStartPre=/sbin/iptables -A INPUT -i vmbr1 -j ACCEPT

# Allow Tailscale interface (if exists)
ExecStart=/bin/sh -c 'iptables -A INPUT -i tailscale0 -j ACCEPT 2>/dev/null || true'

# IPv6 rules
ExecStartPost=/sbin/ip6tables -P INPUT DROP
ExecStartPost=/sbin/ip6tables -P FORWARD ACCEPT
ExecStartPost=/sbin/ip6tables -P OUTPUT ACCEPT
ExecStartPost=/sbin/ip6tables -A INPUT -i lo -j ACCEPT
ExecStartPost=/sbin/ip6tables -A INPUT -m state --state ESTABLISHED,RELATED -j ACCEPT
ExecStartPost=/sbin/ip6tables -A INPUT -i vmbr0 -j ACCEPT
ExecStartPost=/sbin/ip6tables -A INPUT -i vmbr1 -j ACCEPT
ExecStartPost=/bin/sh -c 'ip6tables -A INPUT -i tailscale0 -j ACCEPT 2>/dev/null || true'

# Stop (restore default policies)
ExecStop=/sbin/iptables -P INPUT ACCEPT
ExecStop=/sbin/iptables -F INPUT
ExecStop=/sbin/ip6tables -P INPUT ACCEPT
ExecStop=/sbin/ip6tables -F INPUT

[Install]
WantedBy=multi-user.target
